# easy_boot/ansible_install/templates/sethostname.service.j2
[Unit]
Description=Set hostname to MAC address
Before=salt-minion.service
[Service]
Type=oneshot
ExecStart=/usr/local/bin/sethostname.sh
[Install]
WantedBy=multi-user.target
